# scripts/ai-bridge/testdata/plist/expected.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ai-bridge.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/path/to/ai-bridge</string>
        <string>daemon</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>AI_BRIDGE_CLI</key>
        <string>claude</string>
        <key>AI_BRIDGE_LAUNCHER</key>
        <string>wezterm</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/ai-bridge-daemon.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/ai-bridge-daemon.log</string>
</dict>
</plist>
